# docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=php-fpm
autostart=true
autorestart=true
stderr_logfile=/var/log/php-fpm.err.log  ; Log errors to a file
stdout_logfile=/var/log/php-fpm.out.log   ; Log output to a file
# environment=PHP_FPM_CONFIG=/usr/local/etc/php-fpm.d/www.conf  ; Specify the config file if needed

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx.err.log  ; Log errors to a file
stdout_logfile=/var/log/nginx.out.log   ; Log output to a file
